# roles/aidxp/templates/supervisor.aidxp.conf
[program:uwsgi]
command=uwsgi /etc/uwsgi/apps-enabled/aidxp.ini
user = {{aidxp_webuser}}
group = {aidxp_webuser}}
autostart=true
autorestart=true
stopsignal=INT
